# Install/stuff/1util/time/tzupdate.service
[Unit]
Description=Run tzupdate

[Service]
Type=oneshot
ExecStart=/usr/local/bin/tzupdate
User=root

