# Install/stuff/0important/1pacman/pacmanmirrorupdate.service
[Unit]
Description=Run pacmanmirrorupdate

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pacmanmirrorupdate
User=root

