# scripts/1util/pacman/pacmanmirrorupdate.timer
[Unit]
Description=Timer for running pacmanmirrorupdate

[Timer]
OnBootSec=30min
OnUnitActiveSec=1w

[Install]
WantedBy=timers.target